# inference/server/app/cherrypy.conf
[supervisord]

[program:speech-to-text-server]
command=gunicorn wsgi --config /wav2vec2-server/gunicorn.conf --chdir=/wav2vec2-server
